# deploy/weather-station.service
[Unit]
Description=Weather Station – BME280 + BH1750 → ThingSpeak
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/weather-station
Environment=VIRTUAL_ENV=/home/pi/weather-station/venv
Environment=PATH=/home/pi/weather-station/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/pi/weather-station/venv/bin/python3 weather_station.py
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
